# debian/munin-varnish.postinst
#!/bin/sh

/usr/sbin/invoke-rc.d munin-node restart
